# systemd/airplay-wyse-alsa-policy.service
[Unit]
Description=AirPlay Wyse ALSA Policy Ensure
Documentation=https://github.com/robertmayen/airplay_wyse
Wants=network-online.target
After=network-online.target airplay-wyse-audio-kmods.service
Before=shairport-sync.service

[Service]
Type=oneshot
ExecStart=/usr/local/libexec/airplay_wyse/alsa-policy-ensure
RemainAfterExit=yes
ConditionPathIsExecutable=/usr/local/libexec/airplay_wyse/alsa-policy-ensure

[Install]
WantedBy=multi-user.target
