# systemd/airplay-avahi.service
[Unit]
Description=Airplay wrapper: restart avahi-daemon safely via broker
Documentation=man:systemd.service(5)

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl restart avahi-daemon.service
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/run /var/run
User=root

[Install]
WantedBy=multi-user.target

