# systemd/airplay-shairport.service
[Unit]
Description=Airplay wrapper: restart shairport-sync safely via broker
Documentation=man:systemd.service(5)

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl restart shairport-sync.service
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/run /var/run
User=root

[Install]
WantedBy=multi-user.target

