# systemd/converge-broker.path
[Unit]
Description=Watch /run/airplay/queue to run root broker

[Path]
DirectoryNotEmpty=/run/airplay/queue
Unit=converge-broker.service
MakeDirectory=yes
DirectoryMode=2770

[Install]
WantedBy=multi-user.target

