# systemd/converge.timer
[Unit]
Description=Periodic AirPlay Wyse Converge

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
RandomizedDelaySec=30s
Persistent=true

[Install]
WantedBy=timers.target
