# systemd/reconcile.service
[Unit]
Description=AirPlay Wyse Reconcile (update + converge)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=root
Group=root
WorkingDirectory=/opt/airplay_wyse
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/opt/airplay_wyse/bin/reconcile
# Systemd sandboxing for security
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictSUIDSGID=yes
RestrictRealtime=yes
RestrictNamespaces=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
SystemCallFilter=@system-service
SystemCallFilter=~@privileged @resources @obsolete @debug @mount
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
ReadWritePaths=/var/lib/airplay_wyse /opt/airplay_wyse /etc /run /usr /var/lib/apt /var/cache/apt /var/lib/dpkg /var/log

[Install]
WantedBy=multi-user.target
