# systemd/overrides/nqptp.service.d/override.conf
[Unit]
Wants=network-online.target
After=network-online.target

[Service]
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
CapabilityBoundingSet=
AmbientCapabilities=

