# roles/airplay/templates/airplay-nowplaying.service.j2
[Unit]
Description=AirPlay now-playing metadata reader ({{ airplay_name }})
# Managed by airplay_wyse Ansible.
After=shairport-sync.service
Wants=shairport-sync.service
StartLimitIntervalSec=0

[Service]
Type=simple
User={{ airplay_service_user }}
Group={{ airplay_service_user }}
RuntimeDirectory=airplay
RuntimeDirectoryMode=0755
Environment=AIRPLAY_METADATA_PIPE={{ airplay_metadata_pipe }}
Environment=AIRPLAY_STATE_DIR=/run/airplay
ExecStart=/usr/local/bin/airplay-nowplaying
# Sandbox: reads the metadata pipe and writes only its RuntimeDirectory.
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
# No networking; only local files and the FIFO.
RestrictAddressFamilies=AF_UNIX
SystemCallFilter=@system-service
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
