# roles/airplay/templates/airplay-health.service.j2
[Unit]
Description=AirPlay Wyse health check

[Service]
Type=oneshot
ExecStart=/usr/local/bin/airplay-doctor --check
