# roles/airplay/templates/shairport-sync.service.j2
[Unit]
Description=Shairport Sync ({{ airplay_name }}) - AirPlay 2 Audio Receiver
# Managed by airplay_wyse Ansible. Hardening + ordering live in 10-airplay.conf.
After=network-online.target avahi-daemon.service sound.target alsa-restore.service local-fs.target systemd-udev-settle.service
Wants=network-online.target avahi-daemon.service
# Cold-start tolerance: a USB DAC can take several seconds to enumerate after
# boot, and shairport-sync exits immediately (output_device_error / ENODEV) if
# its raw hw: device is not yet present. Disable the start rate-limiter so the
# unit keeps retrying until the DAC appears, instead of giving up.
StartLimitIntervalSec=0

[Service]
Type=simple
User={{ airplay_service_user }}
Group={{ airplay_service_user }}
SupplementaryGroups=audio
RuntimeDirectory=shairport-sync
RuntimeDirectoryMode=0755
ExecStart=/usr/local/bin/shairport-sync -c /etc/shairport-sync.conf
Restart=on-failure
RestartSec=5
TimeoutStartSec=90
TimeoutStopSec=10

[Install]
WantedBy=multi-user.target
